# service/hypnotoad/run
#!/bin/sh
cd ../..
rm -f var/hypnotoad.lock.*
export MOJO_DAEMON_DEBUG=$(< config/hypnotoad/debug )
exec hypnotoad -f ./web_site </dev/null 2>&1
